# embedded-linux/K26.1/nguyen-hoang-ngoc-hai/session-06/Exercise_2/logger.service
[Unit]
Description=Logger Service Exercise
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/logger
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target